# deploy/launchd/com.sg_mac.gamenews.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sg_mac.gamenews</string>

  <key>UserName</key>
  <string>sg_mac</string>

  <key>WorkingDirectory</key>
  <string>/Users/sg_mac/gamenews</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sg_mac/gamenews/deploy/start-prod.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/sg_mac</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/sg_mac/Library/Logs/com.sg_mac.gamenews.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/sg_mac/Library/Logs/com.sg_mac.gamenews.err.log</string>
</dict>
</plist>
